# ansible/old/files/puma.service
[Unit]
Description=Puma HTTP Server
After=network.target

[Service]
Type=simple
EnvironmentFile=/home/appuser/db_config
User=appuser
<<<<<<< HEAD
WorkingDirectory=/home/appuser/reddit
=======
WorkingDirectory=/tmp/reddit
>>>>>>> ansible-2
ExecStart=/bin/bash -lc 'puma'
Restart=always

[Install]
WantedBy=multi-user.target
